# src/samples/ld.so.preload
/home/vivek/projects/linux-dlp/src/samples/build/bin/libAPIHijack.so
